# docker/files/config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.0]
command=/usr/sbin/php-fpm7.0 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"