# {{cookiecutter.project_name}}/docker/files/config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.2]
command=/usr/sbin/php-fpm7.2 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"